# vacancy-bot.service
[Unit]
Description=Vacancy Bot — агрегатор фрилансерских вакансий
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/vacancy-bot
ExecStart=/opt/vacancy-bot/venv/bin/python -m bot.main
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
